# jjjj.conf
[program:jjjj]
command=python3 /home/curator2/jjjj/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/jjjj/test.err.log
stdout_logfile=/home/curator2/jjjj/test.out.log